# roles/nginx/files/certbot-renewal.timer
[Unit]
Description=Timer for Certbot Renewal

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=multi-user.target
