# roles/nginx/files/certbot-renewal.service
[Unit]
Description=Certbot Renewal

[Service]
ExecStart=certbot renew --post-hook "docker container restart nginx-rev-proxy"
