# src/dns/sources/duckdns.timer
[Unit]
Description=Run DuckDNS DDNS service every 10 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
