# roles/notifications/files/notification-startup.service
[Unit]
Description=Notify Telegram channel, server has started running

[Service]
ExecStart=/usr/local/bin/notification-startup.sh
