# roles/notifications/files/notification-startup.timer
[Unit]
Description=Run start-up notification service on server boot

[Timer]
Unit=notification-startup.service
OnCalendar=*-*-* 8:00:00
Persistent=true

[Install]
WantedBy=timers.target
